# sv/nanoklogd/run
#!/bin/sh
exec chpst -u nobody:nogroup nanoklogd
